# video-splitter-bot.service
[Unit]
Description=Telegram Video Splitter Bot
After=network.target

[Service]
Type=simple
User=root
# Adjust the path to your installation
WorkingDirectory=/opt/video-splitter-bot
# Adjust the path to your python executable (venv)
ExecStart=/opt/video-splitter-bot/venv/bin/python -m bot.main
Restart=always
RestartSec=10
EnvironmentFile=/opt/video-splitter-bot/.env

[Install]
WantedBy=multi-user.target
